# rules/SIGMA-5678-905.yml
# source: test6.yml
id: SIGMA-5678-905
name: Suspicious Scheduled Task Creation or Modification
description: Detects suspicious use of schtasks.exe for creating or modifying scheduled
  tasks, potentially indicating persistence or execution
tactics:
- persistence
- execution
techniques:
- ATTACK.T1053.005
severity: High
query: 'imProcessCreate

  | where ((TargetProcessCommandLine contains "schtasks /create" or TargetProcessCommandLine
  contains "schtasks /change" or TargetProcessCommandLine contains "schtasks /run")
  and TargetProcessName endswith "\\schtasks.exe") and (not((TargetProcessCommandLine
  contains "schtasks /create /tn \"Microsoft\\" or TargetProcessCommandLine contains
  "schtasks /create /tn \"Windows\\")))'
status: test
version: 1.0
clients:
- '*'
tags: []
author: dhif
last_modified: '2025-09-21T19:08:10'
